# debian/package/debian/sample.cron.d
0 9  * * * root /root/asdf.sh > /dev/null 2>&1
*/5 * * * * root /root/asdf.sh > /dev/null 2>&1
